# node/crontab/rotate_mongo_log
0 */12 * * *  mongo admin --eval "printjson(db.runCommand(\"logRotate\"))" --port 27017
