# node/etc/init/server.conf
#!upstart
description "server server"
author      "Andrew Yasinsky"

# used to be: start on startup
# until we found some mounts weren't ready yet while booting:
start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

script
    # Not sure why $HOME is needed, but it is:
    export HOME="/root"

    exec /var/fortihealth/server >> /var/log/server.log 2>&1
end script

pre-start script
    # Date format same as (new Date()).toISOString() for consistency
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Starting" >> /var/log/server.log
end script

pre-stop script
    rm /var/run/server.pid
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Stopping" >> /var/log/server.log
end script
